# crontab
1 0 * * * game-generator
1 12 * * * game-ender && stat-poster
0 1 * * * store_backup.sh
