# crontab
1 0 * * * game-generator
1 12 * * * game-ender && stat-poster
